# backup/button-backup.timer
[Unit]
Description=Run button-backup.service nightly at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
