# scripts/button-deploy.service
[Unit]
Description=Button deployment
After=network.target

[Service]
Type=oneshot
User=button
EnvironmentFile=%h/button.env
ExecStart=%h/button/scripts/deploy.sh --env button
StandardOutput=journal
StandardError=journal
